# var/lib/dpkg/info/binfmt-support.list
/var/cache/binfmts
/var/lib/binfmts
/etc
/etc/init.d
/etc/init.d/binfmt-support
