# squeeze-hindi/var/lib/dpkg/info/apmd.conffiles
/etc/init.d/apmd
/etc/logcheck/ignore.d.paranoid/apmd
/etc/apm/apmd_proxy
/etc/default/apmd
