# squeeze-hindi/var/lib/dpkg/info/apmd.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/apmd" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d apmd stop || exit $?
	else
		/etc/init.d/apmd stop || exit $?
	fi
fi
# End automatically added section
