# squeeze-hindi/var/lib/dpkg/info/apmd.list
/.
/etc
/etc/init.d
/etc/init.d/apmd
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/apmd
/etc/apm
/etc/apm/apmd_proxy
/etc/apm/suspend.d
/etc/apm/other.d
/etc/apm/event.d
/etc/apm/scripts.d
/etc/apm/resume.d
/etc/default
/etc/default/apmd
/usr
/usr/share
/usr/share/apmd
/usr/share/apmd/apmd_proxy.conf
/usr/share/man
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/apmsleep.1.gz
/usr/share/man/man8
/usr/share/man/man8/apmd.8.gz
/usr/share/man/man1
/usr/share/man/man1/apm.1.gz
/usr/share/man/man1/apmsleep.1.gz
/usr/share/doc
/usr/share/doc/apmd
/usr/share/doc/apmd/AUTHORS
/usr/share/doc/apmd/copyright
/usr/share/doc/apmd/apmsleep.README
/usr/share/doc/apmd/changelog.gz
/usr/share/doc/apmd/README
/usr/share/doc/apmd/README.Debian
/usr/share/doc/apmd/changelog.Debian.gz
/usr/bin
/usr/bin/apmsleep
/usr/bin/apm
/usr/sbin
/usr/sbin/apmd
