# squeeze-hindi/var/lib/dpkg/info/gdm.list
/var/lib/gdm
/var/log/gdm
/etc
/etc/pam.d
/etc/pam.d/gdm-autologin
/etc/pam.d/gdm
/etc/init.d
/etc/init.d/gdm
/etc/insserv.conf.d
/etc/insserv.conf.d/gdm
/etc/gdm
/etc/gdm/Xsession
/etc/gdm/gdm.conf
/etc/gdm/XKeepsCrashing
/etc/gdm/modules
/etc/gdm/modules/factory-AccessKeyMouseEvents
/etc/gdm/modules/AccessDwellMouseEvents
/etc/gdm/modules/AccessKeyMouseEvents
/etc/gdm/modules/factory-AccessDwellMouseEvents
/etc/gdm/PreSession
/etc/gdm/PreSession/Default
/etc/gdm/PostSession
/etc/gdm/PostSession/Default
/etc/gdm/PostLogin
/etc/gdm/PostLogin/Default.sample
/etc/gdm/locale.conf
/etc/gdm/Init
/etc/gdm/Init/Default
