# squeeze-hindi/var/lib/dpkg/info/module-init-tools.conffiles
/etc/modprobe.d/aliases.conf
/etc/init.d/module-init-tools
